# hermes.service
[Unit]
Description=Hermes Malware Detection Service
Documentation=https://github.com/yourusername/hermes
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/hermes/app
ExecStart=/usr/bin/python3 /opt/hermes/app/main_linux.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes
ProtectSystem=full
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
